# luci-app-customlogo/root/etc/init.d/customlogo
#!/bin/sh /etc/rc.common

START=99
USE_PROCD=1

BACKUP_SUFFIX=".backup"

backup_file() {
    local file=$1
    if [ -f "$file" ]; then
        if [ ! -f "${file}${BACKUP_SUFFIX}" ]; then
            cp -f "$file" "${file}${BACKUP_SUFFIX}" 2>/dev/null || true
        fi
    fi
}

restore_file() {
    local file=$1
    if [ -f "${file}${BACKUP_SUFFIX}" ]; then
        cp -f "${file}${BACKUP_SUFFIX}" "$file" 2>/dev/null || true
    fi
}

apply_logos() {
    local enable favicon logo ext injected_file css_file

    config_load customlogo
    config_get enable main enable "0"
    config_get favicon main favicon ""
    config_get logo main logo ""

    # 1. 绝对安全的全局还原
    restore_file "/www/favicon.ico"
    restore_file "/www/luci-static/bootstrap/logo.svg"
    restore_file "/www/luci-static/bootstrap/logo_48.png"
    restore_file "/www/luci-static/bootstrap/cascade.css"
    restore_file "/www/luci-static/argon/img/logo.png"
    restore_file "/www/luci-static/argon/img/logo.svg"
    restore_file "/www/luci-static/argon/img/logo-dark.svg"

    # 安全清理 CSS
    css_file="/www/luci-static/bootstrap/cascade.css"
    if [ -f "$css_file" ]; then
        sed -i '/\/\* CUSTOM_LOGO_START \*\//,/\/\* CUSTOM_LOGO_END \*\//d' "$css_file" 2>/dev/null || true
    fi
    rm -f /www/luci-static/bootstrap/custom_injected_logo.* 2>/dev/null || true

    # 2. 只有启用状态才叠加新文件 (解耦 Favicon 和 Logo 的处理逻辑)
    if [ "$enable" = "1" ]; then
        # 确保目录和文件权限
        if [ -d "/etc/customlogo" ]; then
            chmod 644 /etc/customlogo/* 2>/dev/null || true
        fi

        # 独立处理 Favicon
        if [ -n "$favicon" ]; then
            if [ -f "$favicon" ]; then
                backup_file "/www/favicon.ico"
                cp -f "$favicon" "/www/favicon.ico" 2>/dev/null || true
            fi
        fi

        # 独立处理 Logo
        if [ -n "$logo" ]; then
            if [ -f "$logo" ]; then
                ext="${logo##*.}"
                
                # Bootstrap 适配
                if [ -d "/www/luci-static/bootstrap" ]; then
                    backup_file "/www/luci-static/bootstrap/logo.svg"
                    backup_file "/www/luci-static/bootstrap/logo_48.png"
                    backup_file "/www/luci-static/bootstrap/cascade.css"
                    
                    cp -f "$logo" "/www/luci-static/bootstrap/logo.svg" 2>/dev/null || true
                    cp -f "$logo" "/www/luci-static/bootstrap/logo_48.png" 2>/dev/null || true
                    
                    injected_file="custom_injected_logo.$ext"
                    cp -f "$logo" "/www/luci-static/bootstrap/$injected_file" 2>/dev/null || true
                    
                    cat >> "$css_file" <<EOF
/* CUSTOM_LOGO_START */
header .brand {
    font-size: 0 !important;
    color: transparent !important;
    text-shadow: none !important;
    background-image: url('$injected_file') !important;
    background-size: contain !important;
    background-repeat: no-repeat !important;
    background-position: left center !important;
    min-width: 150px !important;
    padding-top: 8px !important;
    padding-bottom: 12px !important;
}
/* CUSTOM_LOGO_END */
EOF
                fi
                
                # Argon 适配
                if [ -d "/www/luci-static/argon" ]; then
                    backup_file "/www/luci-static/argon/img/logo.png"
                    backup_file "/www/luci-static/argon/img/logo.svg"
                    backup_file "/www/luci-static/argon/img/logo-dark.svg"
                    
                    cp -f "$logo" "/www/luci-static/argon/img/logo.png" 2>/dev/null || true
                    cp -f "$logo" "/www/luci-static/argon/img/logo.svg" 2>/dev/null || true
                    cp -f "$logo" "/www/luci-static/argon/img/logo-dark.svg" 2>/dev/null || true
                fi
            fi
        fi
    fi

    # 3. 最终的权限兜底与缓存销毁
    if [ -f "/www/favicon.ico" ]; then
        chmod 644 "/www/favicon.ico" 2>/dev/null || true
    fi
    
    if [ -d "/www/luci-static/argon" ]; then
        chmod 644 /www/luci-static/argon/img/logo* 2>/dev/null || true
    fi
    
    if [ -d "/www/luci-static/bootstrap" ]; then
        chmod 644 /www/luci-static/bootstrap/* 2>/dev/null || true
        rm -f "/www/luci-static/bootstrap/cascade.css.gz" 2>/dev/null || true
        rm -f "/www/luci-static/bootstrap/cascade.css.br" 2>/dev/null || true
    fi
}

start_service() {
    mkdir -p /etc/customlogo
    apply_logos
}

service_triggers() {
    procd_add_reload_trigger "customlogo"
}

reload_service() {
    apply_logos
}
