# luci-app-customlogo/root/etc/init.d/customlogo
#!/bin/sh /etc/rc.common

START=99
USE_PROCD=1

BACKUP_SUFFIX=".backup"

backup_file() {
    local file=$1
    if [ -f "$file" ]; then
        if [ ! -f "${file}${BACKUP_SUFFIX}" ]; then
            cp -f "$file" "${file}${BACKUP_SUFFIX}" 2>/dev/null || true
        fi
    fi
}

restore_file() {
    local file=$1
    if [ -f "${file}${BACKUP_SUFFIX}" ]; then
        cp -f "${file}${BACKUP_SUFFIX}" "$file" 2>/dev/null || true
    fi
}

apply_logos() {
    local enable favicon logo ext injected_file css_file

    config_load customlogo
    config_get enable main enable "0"
    config_get favicon main favicon ""
    config_get logo main logo ""

    # 1. 绝对安全的全局还原 (清除所有主题的修改痕迹)
    restore_file "/www/favicon.ico"
    restore_file "/www/luci-static/bootstrap/logo.svg"
    restore_file "/www/luci-static/bootstrap/logo_48.png"
    restore_file "/www/luci-static/bootstrap/cascade.css"
    restore_file "/www/luci-static/argon/img/logo.png"
    restore_file "/www/luci-static/argon/img/logo.svg"
    restore_file "/www/luci-static/argon/img/logo-dark.svg"

    css_file="/www/luci-static/bootstrap/cascade.css"
    if [ -f "$css_file" ]; then
        sed -i '/\/\* CUSTOM_LOGO_START \*\//,/\/\* CUSTOM_LOGO_END \*\//d' "$css_file" 2>/dev/null || true
    fi
    rm -f /www/luci-static/bootstrap/custom_injected_logo.* 2>/dev/null || true

    # 2. 核心逻辑执行
    if [ "$enable" = "1" ]; then
        
        if [ -d "/etc/customlogo" ]; then
            chmod 644 /etc/customlogo/* 2>/dev/null || true
        fi

        # ==========================================
        # 【完全解耦】专门处理 Favicon
        # ==========================================
        if [ -n "$favicon" ] && [ -f "$favicon" ]; then
            backup_file "/www/favicon.ico"
            cp -f "$favicon" "/www/favicon.ico" 2>/dev/null || true
            
            # 【杀手锏】截胡 Bootstrap 硬编码的 Favicon 路径
            if [ -d "/www/luci-static/bootstrap" ]; then
                backup_file "/www/luci-static/bootstrap/logo.svg"
                backup_file "/www/luci-static/bootstrap/logo_48.png"
                cp -f "$favicon" "/www/luci-static/bootstrap/logo.svg" 2>/dev/null || true
                cp -f "$favicon" "/www/luci-static/bootstrap/logo_48.png" 2>/dev/null || true
            fi
        fi

        # ==========================================
        # 【完全解耦】专门处理 Navigation Bar Logo
        # ==========================================
        if [ -n "$logo" ] && [ -f "$logo" ]; then
            ext="${logo##*.}"
            
            # Bootstrap 主题
            if [ -d "/www/luci-static/bootstrap" ]; then
                # 【注意】我们不再替换 logo.svg 了，而是专属生成注入文件
                backup_file "/www/luci-static/bootstrap/cascade.css"
                
                injected_file="custom_injected_logo.$ext"
                cp -f "$logo" "/www/luci-static/bootstrap/$injected_file" 2>/dev/null || true
                
                local ts=$(date +%s)
                cat >> "$css_file" <<EOF
/* CUSTOM_LOGO_START */
header .brand {
    display: flex !important;
    align-items: center !important;
    margin-left: 0 !important;
    padding-left: 0 !important;
    padding-right: 12px !important;
    padding-top: 0 !important;
    padding-bottom: 0 !important;
    min-width: auto !important;
    width: auto !important;
    height: 40px !important;
    font-size: 20px !important;
    color: #fff !important;
    background-image: none !important;
}
header .brand::before {
    content: url('$injected_file?v=$ts') !important;
    display: block !important;
    height: 28px !important;
    width: auto !important;
    object-fit: contain !important;
    margin-right: 8px !important;
}
/* CUSTOM_LOGO_END */
EOF
            fi
            
            # Argon 主题 (它本身 Navbar 就是图片，保持物理替换)
            if [ -d "/www/luci-static/argon" ]; then
                backup_file "/www/luci-static/argon/img/logo.png"
                backup_file "/www/luci-static/argon/img/logo.svg"
                backup_file "/www/luci-static/argon/img/logo-dark.svg"
                
                cp -f "$logo" "/www/luci-static/argon/img/logo.png" 2>/dev/null || true
                cp -f "$logo" "/www/luci-static/argon/img/logo.svg" 2>/dev/null || true
                cp -f "$logo" "/www/luci-static/argon/img/logo-dark.svg" 2>/dev/null || true
            fi
        fi
    fi

    # 3. 最终的权限兜底与压缩缓存销毁
    [ -f "/www/favicon.ico" ] && chmod 644 "/www/favicon.ico" 2>/dev/null || true
    [ -d "/www/luci-static/argon" ] && chmod 644 /www/luci-static/argon/img/logo* 2>/dev/null || true
    
    if [ -d "/www/luci-static/bootstrap" ]; then
        find /www/luci-static/bootstrap/ -maxdepth 1 -type f -exec chmod 644 {} \; 2>/dev/null || true
        rm -f "/www/luci-static/bootstrap/cascade.css.gz" 2>/dev/null || true
        rm -f "/www/luci-static/bootstrap/cascade.css.br" 2>/dev/null || true
    fi
}

start_service() {
    mkdir -p /etc/customlogo
    apply_logos
}

service_triggers() {
    procd_add_reload_trigger "customlogo"
}

reload_service() {
    apply_logos
}
