# luci-app-customlogo/root/etc/init.d/customlogo
#!/bin/sh /etc/rc.common

START=99
USE_PROCD=1

BACKUP_SUFFIX=".backup"

backup_file() {
    local file=$1
    if [ -f "$file" ] && [ ! -f "${file}${BACKUP_SUFFIX}" ]; then
        cp "$file" "${file}${BACKUP_SUFFIX}"
    fi
}

restore_file() {
    local file=$1
    if [ -f "${file}${BACKUP_SUFFIX}" ]; then
        cp "${file}${BACKUP_SUFFIX}" "$file"
    fi
}

apply_logos() {
    local enable
    local favicon
    local logo

    config_load customlogo
    config_get enable main enable "0"
    config_get favicon main favicon ""
    config_get logo main logo ""

    # 1. 【重构核心】无论当前什么状态，先执行全局“大扫除”，恢复到最纯净的原生状态
    restore_file "/www/favicon.ico"
    restore_file "/www/luci-static/bootstrap/logo.svg"
    restore_file "/www/luci-static/bootstrap/logo_48.png"
    restore_file "/www/luci-static/bootstrap/cascade.css"
    restore_file "/www/luci-static/argon/img/logo.png"
    restore_file "/www/luci-static/argon/img/logo.svg"
    restore_file "/www/luci-static/argon/img/logo-dark.svg"

    # 安全地清理掉任何 CSS 注入残留和之前生成的背景图片
    if [ -f "/www/luci-static/bootstrap/cascade.css" ]; then
        sed -i '/\/\* CUSTOM_LOGO_START \*\//,/\/\* CUSTOM_LOGO_END \*\//d' "/www/luci-static/bootstrap/cascade.css" 2>/dev/null
    fi
    rm -f /www/luci-static/bootstrap/custom_injected_logo.* 2>/dev/null

    # 2. 如果启用了插件，则在纯净的基础上叠加用户的自定义文件
    if [ "$enable" = "1" ]; then
        # 确保 Web 服务器能读取用户上传到 /etc 目录的文件
        chmod 644 /etc/customlogo/* 2>/dev/null

        # 处理 Favicon
        if [ -n "$favicon" ] && [ -f "$favicon" ]; then
            backup_file "/www/favicon.ico"
            cp "$favicon" "/www/favicon.ico" 2>/dev/null
        fi

        # 处理 Logo
        if [ -n "$logo" ] && [ -f "$logo" ]; then
            # 动态获取真实的扩展名 (比如 png, svg, ico)
            local ext="${logo##*.}"
            
            # Bootstrap 主题适配
            if [ -d "/www/luci-static/bootstrap" ]; then
                backup_file "/www/luci-static/bootstrap/logo.svg"
                backup_file "/www/luci-static/bootstrap/logo_48.png"
                backup_file "/www/luci-static/bootstrap/cascade.css"
                
                cp "$logo" "/www/luci-static/bootstrap/logo.svg" 2>/dev/null
                cp "$logo" "/www/luci-static/bootstrap/logo_48.png" 2>/dev/null
                
                # 复制为带真实后缀的专用注入图片
                local injected_file="custom_injected_logo.$ext"
                cp "$logo" "/www/luci-static/bootstrap/$injected_file" 2>/dev/null
                
                local css_file="/www/luci-static/bootstrap/cascade.css"
                cat >> "$css_file" <<EOF
/* CUSTOM_LOGO_START */
header .brand {
    font-size: 0 !important;
    color: transparent !important;
    text-shadow: none !important;
    background-image: url('$injected_file') !important;
    background-size: contain !important;
    background-repeat: no-repeat !important;
    background-position: left center !important;
    min-width: 150px !important;
    padding-top: 8px !important;
    padding-bottom: 12px !important;
}
/* CUSTOM_LOGO_END */
EOF
            fi
            
            # Argon 主题适配
            if [ -d "/www/luci-static/argon" ]; then
                backup_file "/www/luci-static/argon/img/logo.png"
                backup_file "/www/luci-static/argon/img/logo.svg"
                backup_file "/www/luci-static/argon/img/logo-dark.svg"
                
                cp "$logo" "/www/luci-static/argon/img/logo.png" 2>/dev/null
                cp "$logo" "/www/luci-static/argon/img/logo.svg" 2>/dev/null
                cp "$logo" "/www/luci-static/argon/img/logo-dark.svg" 2>/dev/null
            fi
        fi
    fi

    # 3. 【最重要的一步】全局权限护城河与缓存销毁
    [ -f "/www/favicon.ico" ] && chmod 644 "/www/favicon.ico" 2>/dev/null
    [ -d "/www/luci-static/argon" ] && chmod 644 /www/luci-static/argon/img/logo* 2>/dev/null
    
    if [ -f "/www/luci-static/bootstrap/cascade.css" ]; then
        # 强制所有被操作过的文件恢复 644 权限，彻底终结 CSS 裸奔 Bug
        chmod 644 /www/luci-static/bootstrap/* 2>/dev/null
        chmod 644 "/www/luci-static/bootstrap/cascade.css" 2>/dev/null
        
        # 删除预编译的压缩包，强制 Web 服务器重新读取带有你最新设置的 CSS
        rm -f "/www/luci-static/bootstrap/cascade.css.gz" "/www/luci-static/bootstrap/cascade.css.br" 2>/dev/null
    fi
}

start_service() {
    mkdir -p /etc/customlogo
    apply_logos
}

service_triggers() {
    procd_add_reload_trigger "customlogo"
}

reload_service() {
    apply_logos
}
